# systemd/teams-notifications.service
[Unit]
Description=Teams Notifications Tray Daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/Projects/teams-notifications/daemon/.venv/bin/python -m teams_notifications.main
WorkingDirectory=%h/Projects/teams-notifications/daemon
Restart=always
RestartSec=5
Environment=QT_QPA_PLATFORM=xcb
PassEnvironment=DISPLAY XAUTHORITY DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR

[Install]
WantedBy=graphical-session.target
